# supervisor/conf.d/cron.conf
[program:cron]
command=/usr/sbin/cron -f
numprocs=1
directory=/
autostart=true
autorestart=true
